# scripts/systemd/gtos-rpc@.service
[Unit]
Description=GTOS RPC %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tomi
Group=tomi
WorkingDirectory=/data/gtos
EnvironmentFile=/data/gtos/rpc%i/rpc.env
ExecStart=/bin/bash -lc 'exec /usr/local/bin/gtos --config "$GTOS_CONFIG" --verbosity "$GTOS_VERBOSITY" --bootnodes "$GTOS_BOOTNODES" --netrestrict 127.0.0.0/8 --nat none ${GTOS_EXTRA_FLAGS}'
Restart=always
RestartSec=2
LimitNOFILE=1048576
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
