# scripts/systemd/gtos-rpc@.service
[Unit]
Description=GTOS RPC %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tomi
Group=tomi
WorkingDirectory=/data/gtos
EnvironmentFile=/data/gtos/rpc%i/rpc.env
ExecStart=/bin/bash -lc 'exec /usr/local/bin/gtos --config "$GTOS_CONFIG" --datadir "$GTOS_DATADIR" --networkid "$GTOS_NETWORK_ID" --port "$GTOS_P2P_PORT" --netrestrict 127.0.0.0/8 --nat none --http --http.addr 127.0.0.1 --http.port "$GTOS_HTTP_PORT" --http.api "$GTOS_HTTP_API" --ws --ws.addr 127.0.0.1 --ws.port "$GTOS_WS_PORT" --ws.api "$GTOS_WS_API" --authrpc.addr 127.0.0.1 --authrpc.port "$GTOS_AUTHRPC_PORT" --verbosity "$GTOS_VERBOSITY" --bootnodes "$GTOS_BOOTNODES" ${GTOS_EXTRA_FLAGS}'
Restart=always
RestartSec=2
LimitNOFILE=1048576
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
